# fs/etc/services.d/dns-exporter/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid "$PUID:$PGID" /usr/bin/bind_exporter $EXPORTER_ARGS